# Azure-IoT-Connection/azure-iot.service
[Unit]
Description=Azure IoT Hub Connection Service
# Use network.target only so service starts at boot regardless of SIM/network (no network-online.target).
After=network.target nexus-usb0-network.service
Wants=nexus-usb0-network.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /opt/nexuslocate/bin/iot_service.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target